# recipes-edge/edge-resource-manager/edge-resource-manager/edge-resource-manager-watcher.path
[Unit]
Description=Monitor the changes to identity.json file and restart edge-resource-manager

[Path]
PathChanged=EDGE_DATA/edge_gw_config/identity.json
Unit=edge-resource-manager-watcher.service

[Install]
WantedBy=multi-user.target
